# confdev/files/rc.local
$ANSIBLE_VAULT;1.1;AES256
65336264626435393863356637366165396463326364636135613164303032363662366135323733
3435323338366637633861633430333162616465353264310a353638373835306436363037656434
61333136333139316531303064656361306363336265336661656239666536633138663461323835
3337343139326133350a653139666236643732316466346638323363656636373362356266326563
35623236643438663039623330626234386633323164363063316238643631663265616465376566
65326232336138663766316664313039333362633966633361316130376664353562353434336132
37376664366439656630353932303734333935616535633837663237613334316133663762613933
35333963363730646430656630333234643132303639383937346534656132303034313835363135
62653364633232336332363133386538373364653964323130303962306635383363376236323130
38393165313334363131303834333239336536373234626639376365663462343365643532303266
32393336363164346234666439353435373062623766663032306438373266623935636461663032
66316364633935616366396132316464373434393831333132313536306236343864616566356362
36356534313761653732383063336337613935656631666139323130613937653835653333346639
65343034626231663237666665666434316333396130333238336566623139353934313734646565
39383239613466373435303837356661323164393234373032343437323136323436653561663166
61633333336230643766316662333438613032356435326232306439303536643763623837306266
61323830353139346531306262663736353839306533323735386662616634613764666463613432
31303333363863303365363438623639373031343031326530343164386562663835366562343638
36303737373661623865313162613865646536653233383832656262656430646334383030353037
63626131643139333731646639633661343864303937396533396464373537656333323734346161
32653265376631613162636534663534323439303931653761333836383265336133663565323432
38636464313038343062653561336339643735306362383164663064343566643865646261316235
33346536633130663535343937393165303336333330373963663333363038383461613130666135
31613465366134653263303436343536666238323537663936353830363935656333626239663630
31663364633261313831333037356139363734363137363865306364376565343133363931616434
30303936333963646164663962353335653930333132373533643563346635343031616263323264
35663763623862663532396461623961376332636332353931626434663235633633323430643836
33633838646533646637656338613339633637313038386534383662313735353561393731353464
63336262323762316530343365616230366263386637323735643565663134623564316661643365
65383939653439373536306230353334333138353463616662616431613966323933613538363430
38383566656363636536663935376139643832656664393032373237616263356162326537386335
33326238613333616566323239306533636662323861613130613539313161393339623734313338
34383739663261376130653431616266316232356661663134393661343164633931353032623338
31326163343935623365353430373331363763356636376334346263343663326262653632613036
38333537656134633361636138346533376233636532623633643862653364616163373864333432
66393661333537616236663764656531383732663862393938313838333664656132626565623230
32653764386436623431383337653430323832353531333166393830396532646536656238323361
31313961336232666535353938666536643332373866396530313635363965313364623033333939
36643337333936643331316165616436663739666332323634386339663831363463316233316561
30326530363965313463356430633434356534626136353264363839393632376535353936643537
30343937386638663763343735343765623761626361346334383632363536383931646662333037
6563323261396636663964303933663362623731343066373034
